# packages/cups-config/files/cups-setup
#!/bin/sh /etc/rc.common
# CUPS 设置脚本 for THDN-PrintServer
# 启动顺序: 95, 停止顺序: 10

START=95
STOP=10

USE_PROCD=1
PROG=/usr/sbin/cupsd
CONF_FILE=/etc/cups/cupsd.conf
PID_FILE=/var/run/cups/cupsd.pid

start_service() {
    # 创建必要的目录
    mkdir -p /var/log/cups
    mkdir -p /var/run/cups
    mkdir -p /var/spool/cups
    mkdir -p /var/cache/cups
    
    # 设置正确的权限
    chown -R nobody:nogroup /var/log/cups
    chown -R nobody:nogroup /var/spool/cups
    chown -R nobody:nogroup /var/cache/cups
    chmod 755 /var/log/cups
    chmod 755 /var/spool/cups
    chmod 755 /var/cache/cups
    
    # 检查配置文件
    if [ ! -f "$CONF_FILE" ]; then
        echo "CUPS 配置文件不存在: $CONF_FILE"
        return 1
    fi
    
    # 启动 CUPS 服务
    procd_open_instance
    procd_set_param command $PROG -f -c $CONF_FILE
    procd_set_param respawn
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param pidfile $PID_FILE
    procd_close_instance
    
    # 等待服务启动
    sleep 2
    
    # 运行打印机检测脚本
    if [ -x /root/detect-printer.sh ]; then
        /root/detect-printer.sh &
    fi
}

stop_service() {
    # 停止 CUPS 服务
    killall cupsd 2>/dev/null
    
    # 清理临时文件
    rm -f /var/run/cups/cupsd.pid
}

reload_service() {
    # 重新加载 CUPS 配置
    if [ -f "$PID_FILE" ]; then
        kill -HUP $(cat "$PID_FILE") 2>/dev/null
    fi
}

service_triggers() {
    procd_add_reload_trigger "cups"
}

boot() {
    # 延迟启动，确保 USB 子系统准备就绪
    sleep 10
    start
}
